# ru/2026/06/a807265d-2e26-5de0-bdc3-921efb6341b3/victoria-logs.service
[Unit]
Description=VictoriaLogs - High-performance log management and analytics
After=network.target

[Service]
Type=simple
User=victoria-logs
Group=victoria-logs
LimitNOFILE=65536
ExecStart=/opt/victoria-logs/victoria-logs-prod \
  -storageDataPath=/var/lib/victoria-logs \
  -retentionPeriod=30d \
  -syslog.listenAddr.tcp=:29514 \
  -syslog.listenAddr.unix=/run/victoria-logs/victoria-logs.sock
Restart=on-failure

[Install]
WantedBy=multi-user.target
