# ru/2026/03/8adf29e1-b2b4-5714-93a6-31053e3fe67f/3proxy.service
[Unit]
Description=3proxy Proxy Server
After=network.target

[Service]
Type=forking
User=proxy3
Group=proxy3
ExecStart=/usr/local/bin/3proxy /usr/local/etc/3proxy.cfg
ExecReload=/bin/kill -SIGUSR1 $MAINPID
KillMode=process
Restart=on-failure

[Install]
WantedBy=multi-user.target
